# srv.service
[Unit]
Description=Shukarsh Enterprises Product Showcase
After=network.target

[Service]
Type=simple
User=exedev
WorkingDirectory=/home/exedev/shukarsh
Environment=ADMIN_PASSWORD=shukarsh2026
ExecStart=/home/exedev/shukarsh/shukarsh-server --listen :8000
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
